# public/templates/upstart.conf.erb
description '<%= @app.description %>'
author "<%= @app.author || HarbourCrane::Application::NAME  %>"
start on filesystem and started docker
stop on runlevel [!2345]
respawn

env LOGFILE=<%= log_app_file @app.id %>
env CONFIG_DIR=<%= compose_app_dir @app.id %>

pre-start script
  echo "[`date`] <%= @app.name %> Container Starting" >> $LOGFILE
end script

script
  cd $CONFIG_DIR
  COMPOSE_FILE=docker-compose.yml /usr/local/bin/docker-compose up -d
end script

pre-stop script
  echo "[`date`] <%= @app.name %> Container Stopping" >> $LOGFILE
  cd $CONFIG_DIR
  COMPOSE_FILE=docker-compose.yml /usr/local/bin/docker-compose stop
end script

post-stop script
  echo "[`date`] <%= @app.name %> Container Stopped" >> $LOGFILE
end script
